# root/etc/services.d/luadch/run
#!/usr/bin/with-contenv bash
cd /luadch
./luadch -d
